# systemd/claude-revolver-monitor.timer
[Unit]
Description=Poll Claude account usage every 15 minutes
Documentation=https://github.com/ds1sqe/claude-revolver

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
