# systemd/claude-revolver-monitor.timer
[Unit]
Description=Poll Claude account usage periodically
Documentation=https://github.com/ds1sqe/claude-revolver

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
Persistent=true

[Install]
WantedBy=timers.target
